# roles/systemd_units/files/user/notmuch.timer
[Unit]
Description=Index Mail every hour

[Timer]
OnBootSec= 5min
OnUnitActiveSec= 60min
Unit=notmuch.service

[Install]
WantedBy=timers.target
